# tlsdated.service
[Unit]
Description=Secure parasitic rdate replacement

[Service]
ExecStart=/usr/sbin/tlsdated
Type=forking
