# tlsdated.service
[Unit]
Description=Secure parasitic time daemon

[Service]
ExecStart=/usr/sbin/tlsdated
Type=forking
