# tlsdated.service
[Unit]
Description=Secure parasitic rdate replacement

[Service]
ExecStart=/usr/bin/tlsdated
Type=forking
